# arr-appliance/config/arr-disk-guard.service
[Unit]
Description=Run arr disk guard

[Service]
Type=oneshot
ExecStart=/opt/arr-appliance/scripts/guard-disk.sh
