# arr-appliance/config/arr-appliance.service
[Unit]
Description=*arr appliance (Prowlarr, Sonarr, Radarr, Bazarr, qBittorrent)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/arr-appliance
ExecStart=/opt/arr-appliance/scripts/arr-up.sh
ExecStop=/opt/arr-appliance/scripts/arr-down.sh
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
